# .config/systemd/user/kanshi.service
[Unit]
Description=Wayland Display Management Service

[Service]
ExecStart=kanshi

[Install]
WantedBy=sway-session.target

